# scripts/sap-sync-logistica.cron.example
# =============================================================================
# Crontab: sync SAP → SQL Server para KPI_PROD_LOGISTICA
# Fecha inicial: SYNC_FECHA_INICIO en .env.
# =============================================================================
#
#   chmod +x /ruta/al/repo/scripts/sap-sync-logistica-cron.sh
#   crontab -e
#
# Sustituye RUTA_REPO y TU_USUARIO por valores reales.

SHELL=/bin/bash
HOME=/home/TU_USUARIO
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
CRON_TZ=America/Caracas

# Ejemplo: todos los días 10:35 a.m. (hora Caracas)
35 10 * * * /ruta/al/repo/scripts/sap-sync-logistica-cron.sh >> /ruta/al/repo/logs/cron-sync-logistica.log 2>&1
